# fs/etc/supervisor/conf.d/apache.conf
[program:apache2]
priority=5
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
autorestart=unexpected
;startretries=10
;stopwaitsecs=2
;startsecs=2