# fs/etc/supervisor/conf.d/bind.conf
[program:bind9]
priority=6
command=/usr/sbin/named -c /etc/bind/named.conf -u bind -f
autorestart=unexpected
;startretries=10
;stopwaitsecs=2
;startsecs=2