# fs/etc/supervisor/conf.d/mysql.conf
[program:mysql]
priority=4
command=/usr/local/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/sbin/mysqld
autorestart=unexpected
;startretries=10
;stopwaitsecs=2
;startsecs=2