# fs/etc/supervisor/conf.d/ssh.conf
[program:ssh]
priority=3
command=/usr/sbin/sshd -D
autorestart=unexpected
;startretries=10
;stopwaitsecs=2
;startsecs=2