# rootfs/etc/supervisor/conf.d/spamassassin.conf
[program:spamassassin]
command=/usr/sbin/spamd --create-prefs --max-children 5 --helper-home-dir
autostart=false
autorestart=false